# systemd/security-config.service
[Unit]
Description=security-config service

[Service]
Type=simple
ExecStart=/usr/bin/chsmack -r -a "_" /csa

[Install]
WantedBy=multi-user.target

